# roles/common/files/systemd/user/tasksync.service
[Unit]
Description=Sync tasks

[Service]
Type=oneshot
ExecStart=/usr/bin/task sync
