# roles/common/files/systemd/user/tasksync.timer
[Unit]
Description=Sync tasks every 10 minutes
After=network.target

[Timer]
OnBootSec=15sec
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
